# autostart/.config/autostart/kitty.desktop
[Desktop Entry]
Type=Application
Name=Kitty
Exec=kitty
Comment=Load Kitty
