# autostart/.config/autostart/obsidian.desktop
[Desktop Entry]
Type=Application
Name=Obsidian
Exec=flatpak run md.obsidian.Obsidian --background
Comment=Load Obsidian
