# autostart/.config/autostart/protonmail.desktop

[Desktop Entry]
Type=Application
Name=ProtonMail
Exec=proton-mail
Comment=Load ProtonMail
